# crontab
# Barbossa v3.0 - Cron Schedule
# Run every 1 hour to create or review PRs

# Load environment
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

# Every hour at minute 0
# Run as barbossa user (non-root) so Claude CLI works
0 * * * * su - barbossa -c "cd /app && python3 barbossa_simple.py" >> /app/logs/cron.log 2>&1

# Empty line required at end
