# crontab
# Barbossa v3.4 - Cron Schedule
# Tight feedback loop: Tech Lead reviews at :45, Engineer responds at :00
# This keeps work and reviews in sync every hour

# Load environment
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

# Tech Lead - Every hour at minute :45 (REVIEWS FIRST)
# Reviews open PRs, merges good ones, requests changes on others
45 * * * * su - barbossa -c "cd /app && python3 barbossa_tech_lead.py" >> /app/logs/tech_lead_cron.log 2>&1

# Senior Engineer - Every hour at minute :00 (RESPONDS TO FEEDBACK)
# Addresses Tech Lead feedback first, then creates new PRs if queue is clear
0 * * * * su - barbossa -c "cd /app && python3 barbossa_simple.py" >> /app/logs/cron.log 2>&1

# Auditor - Daily at 6:30 AM
# Analyzes system performance, detects patterns, generates recommendations
30 6 * * * su - barbossa -c "cd /app && python3 barbossa_auditor.py --days 7" >> /app/logs/auditor_cron.log 2>&1

# Empty line required at end
